# deploy/systemd/svnportal.service
[Unit]
Description=SvnHub (ASP.NET Core)
After=network.target

[Service]
Type=simple
User=SvnHub
Group=SvnHub
WorkingDirectory=/opt/SvnHub
ExecStart=/usr/bin/dotnet /opt/SvnHub/SvnHub.Web.dll
Restart=always
RestartSec=5
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Environment=ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target

